# mysite/.config/uwsgi/uwsgi.service
[Unit]
Description=uWSGI service
After=syslog.target

[Service]
ExecStart=/home/ubuntu/myvenv/bin/uwsgi -i /home/ubuntu/full-stack/mysite/.config/uwsgi/mysite.ini

Restart=always
KillSignal=SIGQUIT
Type=notify
StandardError=syslog
<<<<<<< HEAD
NotifyAccess=all 
=======
NotifyAccess=all  
>>>>>>> 07ff8f84af9d82d180e9528ee308d9ea2f294b06

[Install]
WantedBy=multi-user.target